# rpm/SOURCES/swgemu-server.service
[Unit]
Description=Star Wars Galaxies server emulator

[Service]
User=swgemu
WorkingDirectory=/opt/swgemu-server/bin/
ExecStart=/opt/swgemu-server/bin/core3
